# scripts/systemd/snowbird-backup.service
[Unit]
Description=Snowbird nightly Postgres backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/snowbird/scripts/pg_backup.sh
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=snowbird-backup
